# docker/server/kubego.service
[Unit]
Description=Kubego watcher

[Service]
Type=simple
Restart=always
RestartSec=1
ExecStart=/app/kubego

[Install]
WantedBy=multi-user.target